# zookeeper/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
user=daemon
priority=0
